# ___lib___systemd___system___reniced.service
[Unit]
Description=reniced

[Service]
ExecStart=/bin/bash -c 'while true; do /usr/bin/reniced; sleep 10; done'

[Install]
WantedBy=multi-user.target
