# copy/__etc__systemd__system__reniced.service
[Unit]
Description=reniced

[Service]
ExecStart=/etc/.phdconf/bin/reniced.sh

[Install]
WantedBy=multi-user.target
